# render-api/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:bridge]
command=node /app/claude-bridge/dist/server.js
directory=/tmp/claude-bridge/workdir
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=TERM
stopwaitsecs=60
environment=NODE_ENV="production"

[program:api]
command=node /app/dist/index.js
directory=/app
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=TERM
stopwaitsecs=30
environment=NODE_ENV="production"
